# deployment/test/simpleapi.ini
[program:simpleapi1]
environment=SIMPLEAPI_CFG="test",SIMPLEAPI_NODE="1"
directory=/data/workspace/simpleapi/
command=/data/workspace/simpleapi/venv/bin/gunicorn -c deployment/test/gunicorn.py app:app
stdout_logfile=/tmp/sup.simpleapi.stdout.log
redirect_stderr=true
autostart=true
autorestart=true

[program:simpleapi2]
environment=SIMPLEAPI_CFG="test",SIMPLEAPI_NODE="2"
directory=/data/workspace/simpleapi/
command=/data/workspace/simpleapi/venv/bin/gunicorn -c deployment/test/gunicorn.py app:app
stdout_logfile=/tmp/sup.simpleapi.stdout.log
redirect_stderr=true
autostart=true
autorestart=true

[program:simpleapi3]
environment=SIMPLEAPI_CFG="test",SIMPLEAPI_NODE="3"
directory=/data/workspace/simpleapi/
command=/data/workspace/simpleapi/venv/bin/gunicorn -c deployment/test/gunicorn.py app:app
stdout_logfile=/tmp/sup.simpleapi.stdout.log
redirect_stderr=true
autostart=true
autorestart=true
